# elasticsearch/install
#!/bin/bash
set -eu

sudo apt-get update -y
sudo apt-get install -y apt-transport-https ca-certificates wget gpg time jq bc python3 python3-pip

# Elastic signing key + repo (idempotent if already present).
if [ ! -f /usr/share/keyrings/elasticsearch-keyring.gpg ]; then
    wget -qO - https://artifacts.elastic.co/GPG-KEY-elasticsearch \
        | sudo gpg --dearmor -o /usr/share/keyrings/elasticsearch-keyring.gpg
fi

if [ ! -f /etc/apt/sources.list.d/elastic-9.x.list ]; then
    echo "deb [signed-by=/usr/share/keyrings/elasticsearch-keyring.gpg] https://artifacts.elastic.co/packages/9.x/apt stable main" \
        | sudo tee /etc/apt/sources.list.d/elastic-9.x.list
    sudo apt-get update -y
fi

sudo apt-get install -y elasticsearch

# load.py uses requests.
pip3 install --quiet --break-system-packages requests || pip3 install --quiet requests

# Disable security; bind on all interfaces; run as a single-node cluster.
#
# The original config set `cluster.initial_master_nodes: ["clickbench"]`
# without a matching `node.name: clickbench`. ES 9.0–9.3 tolerated the
# mismatch (the bootstrap check would still pick the first master-eligible
# node). 9.4 (released 2026-05-05) tightened bootstrap so the listed
# initial master must match an existing `node.name`; otherwise the
# cluster never elects a master, status stays red, and bulk writes get
# rejected with HTTP 503 forever (we saw 200+ batches of 503 followed
# by ES exiting and "Connection refused").
#
# `discovery.type: single-node` is the documented way to skip cluster
# bootstrap entirely for a single-process install — equivalent to what
# we want and resilient to whatever node.name systemd gives the node.
cat <<EOF | sudo tee /etc/elasticsearch/elasticsearch.yml >/dev/null
path.data: /var/lib/elasticsearch
path.logs: /var/log/elasticsearch

xpack.security.enabled: false
xpack.security.http.ssl.enabled: false
xpack.security.transport.ssl.enabled: false

discovery.type: single-node
http.host: 0.0.0.0
EOF

sudo /bin/systemctl daemon-reload
sudo systemctl enable elasticsearch.service
sudo systemctl restart elasticsearch.service
