# playground/clickbench-playground.service
[Unit]
Description=ClickBench Playground API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ClickBench
EnvironmentFile=-/home/ubuntu/ClickBench/playground/.env
ExecStart=/usr/bin/python3 -m playground.server.main
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
